# inst/server/plumber.service
[Unit]
Description=Plumber API

[Service]
ExecStart=/usr/bin/Rscript -e "pr <- plumber::plumb('/var/plumber$PATH$/plumber.R'); $PREFLIGHT$ pr$run(port=$PORT$, swagger=$SWAGGER$)"
Restart=on-abnormal
WorkingDirectory=/var/plumber/$PATH$/

[Install]
WantedBy=multi-user.target
